# meta-lumotive/recipes-bsp/hwrevision/hwrevision/hwrevision.service
[Unit]
Description=Set hwrevision file

[Service]
Type=oneshot
ExecStart=/usr/bin/hwrevision.sh

[Install]
WantedBy=multi-user.target
